# OpenWrt/官方OpenWrt（带固件编译定制流程）/crontabs/root
00 6 * * * reboot



*/2 * * * * /usr/sbin/autologin 学号 密码 天翼账号 天翼密码
